# resources/app.service.in
[Unit]
Description=Spacelink Daemon
After=multi-user.target
Requires=network.target
ConditionVirtualization=!container

[Service]
Type=dbus
BusName=@DBUS_NAME@
ExecStart=@DAEMON_BINARY@
StandardError=null
Restart=on-abort
ProtectSystem=true
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
User=root

[Install]
WantedBy=multi-user.target
